# steam-powerbuttond-chimeraos.service
[Unit]
Description=Steam Power Button Service for ChimeraOS
After=graphical-session.target

[Service]
ExecStart=/home/gamer/.local/bin/steam-powerbuttond

[Install]
WantedBy=multi-user.target
